# etc/supervisord.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autorestart=true
stopsignal=KILL
environment=LD_PRELOAD=""
